# deploy/openclaw-status.timer
# Runs the pipeline shortly after boot and every 6h thereafter. The app polls
# npm/GitHub for new OpenClaw releases, so a few-hour cadence keeps the verdict
# fresh without burning API budget. Adjust OnUnitActiveSec to taste.
[Unit]
Description=Run the OpenClaw Status pipeline periodically

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
Persistent=true

[Install]
WantedBy=timers.target
