# deploy/openclaw-status.service
# Pulls the latest code, then runs the adaptive scheduler *tick*: it decides whether a
# full collect→assess→render assessment is due and runs one only when warranted (a new
# release, or the cadence in config.ASSESS_CADENCE_TIERS). Most ticks are a cheap GitHub
# check with no LLM cost. Triggered hourly by openclaw-status.timer (not enabled alone).
[Unit]
Description=OpenClaw Status — pull latest code + run the scheduler tick
After=network-online.target
Wants=network-online.target
# Hard failure (git/collector crash, LLM timeout) → ping the alert webhook. The deploy
# guard already keeps the last good page live; this is the only out-of-band signal it died.
OnFailure=openclaw-status-failure.service

[Service]
Type=oneshot
User=openclaw
Group=openclaw
WorkingDirectory=/opt/openclaw_status_app
Environment=PYTHONUNBUFFERED=1
# Pull latest code via a wrapper that tolerates a transient offline pull but SURFACES a
# persistent non-fast-forward divergence to the alert webhook — otherwise the box would run
# stale code silently (D30). The "-" stays as a second belt so nothing here blocks the tick.
ExecStartPre=-/opt/openclaw_status_app/deploy/git-pull.sh
ExecStart=/opt/openclaw_status_app/.venv/bin/python run.py tick
# Most ticks finish in seconds; when one decides to run, the full pipeline takes a few
# minutes. Must stay ABOVE config.COLLECT_TIMEOUT_S + PIPELINE_BUDGET_S + render margin
# (480 + 900 + ~60 ≈ 1440) so the in-process budgets degrade gracefully — publishing the
# primary / keeping the last good page — before this SIGKILL can pre-empt them.
TimeoutStartSec=1800
Nice=10
# ── Sandboxing (contain the blast radius) ──
# This unit parses untrusted LLM JSON + GitHub issue/release text and pulls PyPI deps, while
# the openclaw uid holds the secrets (.env) and writes the public web/ tree — so an RCE here
# could exfiltrate keys AND deface the page. Lock it down: no privilege escalation, the rest of
# the filesystem read-only, homes read-only, a private /tmp. Writes are confined to the app dir
# (git pull rewrites the repo; the pipeline writes data/ + web/ — all under it).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/opt/openclaw_status_app
